# config/s6-rc.d/init-script/up
#!/bin/execlineb -P
/bin/ash -c "/init-script.sh"