# config/s6-rc.d/php-fpm/run
#!/bin/execlineb -P
php-fpm{PHP_VERSION} -F -R
